# compile/systemd/civwiki-inflight.timer
[Unit]
Description=Refresh the Civilization Wiki live in-flight overlay every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
