# compile/systemd/wiki-autodeploy.service
[Unit]
Description=Wiki auto-deploy poller (one tick)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /Transpara/transpara-ai/repos/wiki/compile/autodeploy.py
