# compile/systemd/wiki-inflight.service
[Unit]
Description=Wiki live in-flight overlay refresh (one tick)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /Transpara/transpara-ai/repos/wiki/compile/inflight.py
